# server/stories.service
[Unit]
Description=Test Service

[Service]
ExecStart=[path-to]/start.sh

[Install]
WantedBy=multi-user.target